# aws_ork/data/systemd/aws_ork.service
[Unit]
Description=AWS Ork
After=network.target

[Service]
Type=simple
NotifyAccess=all
ExecStart=/usr/local/bin/aws_ork
KillMode=process

[Install]
WantedBy=multi-user.target
